# ma_alert-swiftdialog/payload/Library/LaunchAgents/com.company.migration-alert.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.company.migration-alert</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/CPE/bin/migration-alert</string>
	</array>
	<key>StandardOutPath</key>
	<string>/Library/CPE/Logs/migration-alert.log</string>
	<key>StandardErrorPath</key>
	<string>/Library/CPE/Logs/migration-alert.log</string>
	<key>RunAtLoad</key>
	<true />
	<key>KeepAlive</key>
	<true />
</dict>
</plist>
